# scripts/tiffany-bot.service
[Unit]
Description=Tiffany Bot (Discord)
After=network-online.target warp-svc.service
Wants=network-online.target
# WARP provides the SOCKS5 proxy (127.0.0.1:40000) that yt-dlp needs for YouTube.
# Not a hard Requires= so the bot still starts (news/offers work) if WARP is down.
Wants=warp-svc.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tiffany-bot
ExecStartPre=/bin/bash /opt/tiffany-bot/scripts/kill-orphans.sh --pre-start
# run.sh prefers the venv (Python 3.11+) and falls back to system python3.
ExecStart=/bin/bash /opt/tiffany-bot/scripts/run.sh
Restart=on-failure
RestartSec=10

# Mata TODOS os processos filhos ao parar (resolve processos órfãos)
KillMode=control-group
KillSignal=SIGTERM
FinalKillSignal=SIGKILL
TimeoutStopSec=10
SendSIGKILL=yes

# Variáveis de ambiente
Environment=PYTHONUNBUFFERED=1

# Logs via journalctl (substitui bot.log + nohup)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tiffany-bot

[Install]
WantedBy=multi-user.target
