# scripts/tiffany-bot.service
[Unit]
Description=Tiffany Bot (Discord)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tiffany-bot
ExecStartPre=/bin/bash -c '/usr/bin/pkill -9 -f "/opt/tiffany-bot/launcher.py" 2>/dev/null; /usr/bin/pkill -9 -f "/opt/tiffany-bot/notices.py" 2>/dev/null; /usr/bin/pkill -9 -f "/opt/tiffany-bot/offers.py" 2>/dev/null; sleep 1; rm -f /tmp/tiffany_launcher.lock; exit 0'
ExecStart=/usr/bin/python3 -u launcher.py
Restart=on-failure
RestartSec=10

# Mata TODOS os processos filhos ao parar (resolve processos órfãos)
KillMode=control-group
KillSignal=SIGKILL
TimeoutStopSec=15

# Variáveis de ambiente
Environment=PYTHONUNBUFFERED=1

# Logs via journalctl (substitui bot.log + nohup)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tiffany-bot

[Install]
WantedBy=multi-user.target
